# etc/runit/sv/vmwaretools/run
#!/bin/sh
/etc/init.d/vmware-tools start >/dev/null
exec chpst -b vmwaretools pause
