# etc/runit/sv/vmwaretools/run
#!/bin/sh
/etc/init.d/vmware-tools start
exec chpst -b vmwaretools pause
